# docker/certbot/crontab
42 6 * * * /usr/local/bin/be-certbot renew
